# watchtower-sentinel.service
[Unit]
Description=Watchtower Sentinel
After=network.target

[Service]
WorkingDirectory=/opt/watchtower-sentinel
ExecStart=/opt/watchtower-sentinel/watchtower-sentinel
EnvironmentFile=/opt/watchtower-sentinel/.env
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
